# archive/windowsir/posts/0190-root-cause-analysis.json
{
  "url": "https://windowsir.blogspot.com/2022/04/root-cause-analysis.html",
  "title": "Root Cause Analysis",
  "published": "2022-04-29T16:54:00.002-05:00",
  "updated": "2022-04-29T16:54:54.490-05:00",
  "labels": [],
  "text": "One of the challenges within DFIR, particularly as we've moved to an enterprise approach by leveraging EDR telemetry, is the root cause analysis, or \"RCA\". In short, the challenge is observing malicious activity and determining the root cause; the challenge itself stems from the fact that EDR telemetry is only partial visibility, or that correlating observed malicious activity with causal data not evident or available via EDR telemetry requires additional context, and by extension, additional effort/expenditure of resources. It also requires an additional \"leveling up\" of skillsets.\n Yes, many organizations that deploy EDR tooling also include a means for extracting additional files/data from the endpoint, and what to collect isn't usually in question. Rather, how to truly exploit the collected data is the issue, with the exploitation of that data being the \"levelling up\".\n When malicious activity is observed, or the impact of malicious activity is discovered (via threat hunting, etc.), the challenge then becomes, how do we determine the root cause? In the past, when someone has shared suspicious activity with me and sought guidance as to next steps, I've often asked, \"...did a reboot or login occur just prior to the activity?\" as a means of developing context around the observed activity. Was the executable launched via an auto-start mechanism, such as a Windows service, entry in a Run key, or as a Scheduled Task? The thought process has been to seek causal events for the observed activity, which can be important to not only determine the root cause, but perhaps to also identify a shift in TTPs , once we're to a point where we can arrive at attribution.\n There are a lot of different ways for threat actors to persist on Windows systems, some of which were mentioned earlier in this post. Each has their advantages and disadvantages, and both require and provide different levels of access. For example, creating an entry in the user's Run key means that the malware won't start again until the user logs in, and then runs within the user context. However, if the threat actor has Admin privileges, they can create a Windows service or Scheduled Task, which then will run malware with SYSTEM level privileges. As a result, the persistence mechanism used provides a great deal of context beyond just the entry for the malware.\n To provide further insight into this topic, Krz wrote up an excellent blog post not long ago about how a laptop going on battery (unplugged from the power cord) can impact an investigation.\n Based on the MITRE \" Event Triggered Execution: Screensaver \" page, I created a RegRipper plugin for my local repo to parse the screensaver settings out of the NTUSER.DAT hive, and ran it against a hive extracted from my local system, the output of which follows:",
  "html": "<p style=\"text-align: left;\">One of the challenges within DFIR, particularly as we've moved to an enterprise approach by leveraging EDR telemetry, is the root cause analysis, or \"RCA\". In short, the challenge is observing malicious activity and determining the root cause; the challenge itself stems from the fact that EDR telemetry is only partial visibility, or that correlating observed malicious activity with causal data not evident or available via EDR telemetry requires additional context, and by extension, additional effort/expenditure of resources. It also requires an additional \"leveling up\" of skillsets.&nbsp;</p><p style=\"text-align: left;\">Yes, many organizations that deploy EDR tooling also include a means for extracting additional files/data from the endpoint, and <i>what</i> to collect isn't usually in question. Rather, how to truly exploit the collected data is the issue, with the exploitation of that data being the \"levelling up\".</p><p style=\"text-align: left;\">When malicious activity is observed, or the impact of malicious activity is discovered (via threat hunting, etc.), the challenge then becomes, how do we determine the root cause? In the past, when someone has shared suspicious activity with me and sought guidance as to next steps, I've often asked, \"...did a reboot or login occur just prior to the activity?\" as a means of developing context around the observed activity. Was the executable launched via an auto-start mechanism, such as a Windows service, entry in a <i>Run</i> key, or as a Scheduled Task? The thought process has been to seek causal events for the observed activity, which can be important to not only determine the root cause, but perhaps to also identify a <a href=\"https://www.bleepingcomputer.com/news/security/emotet-malware-now-installs-via-powershell-in-windows-shortcut-files/\">shift in TTPs</a>, once we're to a point where we can arrive at attribution.</p><p style=\"text-align: left;\">There are a lot of different ways for threat actors to persist on Windows systems, some of which were mentioned earlier in this post. Each has their advantages and disadvantages, and both require and provide different levels of access. For example, creating an entry in the user's <i>Run</i> key means that the malware won't start again until the user logs in, and then runs within the user context. However, if the threat actor has Admin privileges, they can create a Windows service or Scheduled Task, which then will run malware with SYSTEM level privileges. As a result, the persistence mechanism used provides a great deal of context beyond just the entry for the malware.&nbsp;</p><p style=\"text-align: left;\">To provide further insight into this topic, <a href=\"https://cyberdefnerd.com/\">Krz</a> wrote up an excellent blog post not long ago about how a <a href=\"https://cyberdefnerd.com/2022/03/27/why-do-the-battery-use-and-the-battery-level-matter-during-the-investigation/\">laptop going on battery</a> (unplugged from the power cord) can impact an investigation.</p><p style=\"text-align: left;\">Based on the MITRE \"<a href=\"https://attack.mitre.org/techniques/T1546/002/\">Event Triggered Execution: Screensaver</a>\" page, I created a RegRipper plugin for my local repo to parse the screensaver settings out of the NTUSER.DAT hive, and ran it against a hive extracted from my local system, the output of which follows:</p><div style=\"text-align: left;\"><span style=\"font-family: courier;\">Launching screensaver v.20220427<br />screensaver v.20220427<br />(NTUSER.DAT) Gets user's screensaver settings<br />MITRE: <a href=\"https://attack.mitre.org/techniques/T1546/002/\">T1546.002</a> (persistence)<br /><br /></span></div><div style=\"text-align: left;\"><span style=\"font-family: courier;\">Control Panel\\Desktop<br />LastWrite: 2021-11-29 11:53:45Z<br /><br /></span></div><div style=\"text-align: left;\"><span style=\"font-family: courier;\">Screensaver is active.<br />SCRNSAVE.exe value not found.<br /><br /></span></div><div style=\"text-align: left;\"><span style=\"font-family: courier;\">Analysis Tip: Threat actors have been observed using the screen saver as a persistent mechanism.<br /><br /></span></div><div style=\"text-align: left;\"><span style=\"font-family: courier;\">Ref: https://cocomelonc.github.io/tutorial/2022/04/26/malware-pers-2.html</span></div><div style=\"text-align: left;\"><span style=\"font-family: courier;\"><br /></span></div><div style=\"text-align: left;\"><span style=\"font-family: inherit;\">Now, if the \"SCRNSAVE.exe\" value (listed above as \"not found\") was set to something that it shouldn't be, we could then use the <i>Desktop</i> key LastWrite time as a pivot point for further analysis. As the <i>Desktop</i> value contains a number of values, we may not be able to specifically attribute the addition or alteration of the SCRNSAVE.exe value to the modification of the key LastWrite time, but it does provide us with some information we can pivot on and leverage in our analysis.&nbsp;</span></div><div style=\"text-align: left;\"><span style=\"font-family: inherit;\"><br /></span></div><div style=\"text-align: left;\"><span style=\"font-family: inherit;\">Another instance to make clear why this example is so compelling...see <a href=\"https://twitter.com/jorgeorchilles/status/1519356926404050945\">Jorge's tweet here</a>, which points us to <a href=\"https://twitter.com/Wietze/status/1519260918236000257\">Wietze's tweet</a>, which in turn points us to <a href=\"https://lolbas-project.github.io/lolbas/Libraries/Desk/\">a new #LOLBAS</a>. Kind of a circuitous route (also apparently available as far back as 2001, from <a href=\"https://twitter.com/pegabizu/status/1519633040473038848\">here</a>), I know, but ultimately what this leads us to is the use of a Control Panel applet (desktop.cpl) to proxy the execution of arbitrary commands. In the example, 'calc.exe' is renamed to 'calc.scr' and launched by calling the \"InstallScreenSaver\" function of the applet. As a result, this persistence mechanism can be established via <i>reg.exe</i> or <i>rundll32.exe</i>; different command lines with different impacts on the system (re: artifacts associated with their execution) but both achieve the same result.</span></div><div style=\"text-align: left;\"><span style=\"font-family: inherit;\"><br /></span></div><div style=\"text-align: left;\"><span style=\"font-family: inherit;\">Further, we now have insight into the level of access achieved by the threat actor, the context under which the malware will run, and when the malware is expected to run. This also provides insights into the <a href=\"https://windowsir.blogspot.com/2022/01/registry-analysis-why.html\">sophistication and intent</a> of the threat actor, and can be applied toward attribution.</span></div>"
}
